# EnergySuper/EnergySuper.service
[Unit]
Description=Energy Systems Supervisor

[Service]
Type=notify
ExecStart=/usr/sbin/EnergySuper/EnergySuper

[Install]
WantedBy=multi-user.target
